# hobo-text/deploy/systemd/hobo-text.service
[Unit]
Description=HoboText + HoboLogo — Text Generation & Logo Tools
After=network.target

[Service]
Type=simple
User=ubuntu
WorkingDirectory=/opt/hobo-text
Environment=NODE_ENV=production
ExecStart=/home/ubuntu/.nvm/versions/node/v20.20.1/bin/node /opt/hobo-text/server/index.js
Restart=always
RestartSec=5
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ProtectHome=false
ReadWritePaths=/opt/hobo-text

[Install]
WantedBy=multi-user.target
