# hobo-tools/deploy/systemd/hobo-tools.service
[Unit]
Description=HoboTools — Central Hub & SSO
After=network.target

[Service]
Type=simple
User=deck
WorkingDirectory=/opt/hobo-tools
Environment=NODE_ENV=production
ExecStart=/usr/bin/node /opt/hobo-tools/server/index.js
Restart=always
RestartSec=5
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ProtectHome=true
ReadWritePaths=/opt/hobo-tools/data

[Install]
WantedBy=multi-user.target
