# hobo-img/deploy/systemd/hobo-img.service
[Unit]
Description=HoboImg — Unified Image Converter & Processing Tools
After=network.target

[Service]
Type=simple
User=ubuntu
WorkingDirectory=/opt/hobo-img
Environment=NODE_ENV=production
ExecStart=/home/ubuntu/.nvm/versions/node/v20.20.1/bin/node /opt/hobo-img/server/index.js
Restart=always
RestartSec=5
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ProtectHome=false
ReadWritePaths=/opt/hobo-img

[Install]
WantedBy=multi-user.target
